# infra/deploy/audit-verify.timer
[Unit]
Description=Hourly AETHER audit chain verification
After=network.target

[Timer]
OnCalendar=hourly
RandomizedDelaySec=120
Persistent=true

[Install]
WantedBy=timers.target
